# src/ibm/hypervisor-network-mgr-src/xyz.openbmc_project.Network.Hypervisor.service.in
[Unit]
Description=Hypervisor Network Manager
Wants=xyz.openbmc_project.biosconfig_manager.service
After=xyz.openbmc_project.biosconfig_manager.service
Wants=pldmd.service
After=pldmd.service

[Service]
ExecStart=/usr/bin/hyp-network-manager
Restart=always
Type=dbus
BusName=xyz.openbmc_project.Network.Hypervisor

[Install]
WantedBy=multi-user.target
