# src/lldp/xyz.openbmc_project.LLDP.service.in
[Unit]
Description=LLDP for Neighbor Discovery
After=lldpd.service
After=network.target

[Service]
ExecStart=/usr/bin/phosphor-lldpd
Restart=always
Type=dbus
BusName=@LLDP_DEFAULT_NETWORK_BUSNAME@

[Install]
WantedBy=multi-user.target
